# build/graphite-web/supervisor/graphite.conf
[program:graphite-web]
command = /usr/sbin/apache2ctl -D FOREGROUND
autorestart = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0

[program:restartapache2]
command = /bin/bash -c "sleep 3 && service apache2 reload"
startsecs = 0
autorestart = false
