# build/elasticsearch-curator/supervisor/curator.conf
[program:cron]
command = cron -f
autorestart = true

[program:tail]
command = tail -f /var/log/curator
autorestart = true