# build/graphite/supervisor/graphite.conf
[program:initdb]
command = /bin/bash -c "chown -R _graphite:_graphite /var/lib/graphite && if [ ! -f /var/lib/graphite/graphite.db ]; then sudo -u _graphite graphite-manage syncdb --noinput; fi"
startsecs = 0
autorestart = false

[program:carbon]
command = /usr/bin/carbon-cache --config=/etc/carbon/carbon.conf --debug start
autorestart = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0

[program:graphite-web]
command = /usr/sbin/apache2ctl -D FOREGROUND
autorestart = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0

[program:restartapache2]
command = /bin/bash -c "sleep 5 && service apache2 reload"
startsecs = 0
autorestart = false
