# build/sensu/supervisor/sensu.conf
[program:rabbitmq]
command = /usr/sbin/rabbitmq-server
autorestart = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0

[program:redis]
command = /usr/bin/redis-server
autorestart = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0

[program:sensu-server]
command = /opt/sensu/bin/sensu-server -c /etc/sensu/config.json -d /etc/sensu -e /etc/sensu/extensions -v -l /var/log/sensu/server.log
autorestart = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0

[program:sensu-api]
command = /opt/sensu/bin/sensu-api -c /etc/sensu/config.json -d /etc/sensu -e /etc/sensu/extensions -v -l /var/log/sensu/api.log
autorestart = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0