# build/carbon/supervisor/carbon.conf
[program:carbon-cache]
command = /usr/bin/carbon-cache --config=/etc/carbon/carbon.conf --debug start
autorestart = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
