# build/sensu-client/supervisor/sensu.conf
[program:sensu-client]
command = /opt/sensu/bin/sensu-client -c /etc/sensu/config.json -d /etc/sensu/conf.d -e /etc/sensu/extensions -v
autorestart = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
environment = PATH=/etc/sensu/plugins:/etc/sensu/handlers:%(ENV_PATH)s