# build/grafana/supervisor/grafana.conf
[program:initconf]
command = cp -f /config/config.js /usr/share/grafana/config.js
startsecs = 0
autorestart = false

[program:grafana]
command = /usr/sbin/apache2ctl -D FOREGROUND
autorestart = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0